# templates/default/cpuminer.upstart.conf.erb
#!upstart
description "cpuminer deamon"

start on network
stop on shutdown

respawn
nice 19
setuid nobody
setgid nogroup

exec /usr/bin/env minerd -o <%= @node['miner']['url'] %> -u <%= @node['miner']['user'] %> -p <%= @node['miner']['password'] %>